# root/etc/services.d/lidarr/run
#!/command/with-contenv
umask "${UMASK}"
exec s6-setuidgid agpsn "${APP_DIR}/lidarr/bin/Lidarr" -nobrowser -data="${CONFIG_DIR}"
